# services/gl-dashboard.service
[Unit]
Description=Ghostlink Dashboard
After=gl-network.service
Requires=gl-network.service

[Service]
Type=simple
User=ghostlink
Group=ghostlink
WorkingDirectory=/opt/ghostlink/dashboard
Environment=GHOSTLINK_CONF=/etc/ghostlink/ghostlink.conf
ExecStart=/opt/ghostlink/venv/bin/gunicorn \
    --bind unix:/run/ghostlink/dashboard.sock \
    --workers 2 \
    --worker-class geventwebsocket.gunicorn.workers.GeventWebSocketWorker \
    --certfile /etc/ghostlink/dashboard.crt \
    --keyfile  /etc/ghostlink/dashboard.key \
    app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=ghostlink.target
