# services/gl-identity.service
[Unit]
Description=Ghostlink Identity Engine
After=gl-system.service network-pre.target
Before=gl-network.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/ghostlink/identity/engine.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=ghostlink.target
